# installers/linux/usr/lib/systemd/user/ai-orchestrator.service
[Unit]
Description=AI Orchestrator — Local AI agent web server
Documentation=https://github.com/pop31-ai/ai-orchestrator
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
ExecStart=%h/.local/share/ai-orchestrator/venv/bin/python -m ai_orchestrator serve
WorkingDirectory=%h/.local/share/ai-orchestrator
Restart=on-failure
RestartSec=5
Environment="PYTHONUNBUFFERED=1"

# Безопасность
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=true

[Install]
WantedBy=default.target